# simple-java-maven-app-master/supervisord.conf
nodaemon=true

[program:my-app]
command=java -jar target/my-app-1.0-SNAPSHOT.jar
directory=/app
autostart=true
autorestart=true